# init/gvb_server.ini
[program: gvb_server]
command=/opt/gvb/server/main
directory=/opt/gvb/server/
autorestart=true ; 程序意外退出是否自动重启
autostart=true ; 是否自动启动
stderr_logfile=/opt/gvb/server/err.log ; 错误日志文件
stdout_logfile=/opt/gvb/server/out.log ; 输出日志文件
user=root ; 进程执行的用户身份
stopsignal=INT
startsecs=1 ; 自动重启间隔
stopasgroup=true ;默认为false,进程被杀死时，是否向这个进程组发送stop信号，包括子进程
killasgroup=true ;默认为false，向进程组发送kill信号，包括子进程